# scripts/com.meta-scheduler.supervisor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.meta-scheduler.supervisor</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>-lc</string>
    <string>cd __REPO_DIR__ &amp;&amp; ./node_modules/.bin/tsx src/index.ts</string>
  </array>

  <key>EnvironmentVariables</key>
  <dict>
    <key>VAULT_DIR</key>
    <string>__REPO_DIR__/private/agents</string>
    <key>BIND_HOST</key>
    <string>127.0.0.1</string>
    <key>PORT</key>
    <string>7421</string>
    <key>RESTART_ENABLED</key>
    <string>false</string>
  </dict>

  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>

  <key>StandardOutPath</key>
  <string>__HOME__/.meta-scheduler/supervisor.out.log</string>
  <key>StandardErrorPath</key>
  <string>__HOME__/.meta-scheduler/supervisor.err.log</string>
</dict>
</plist>
